# examples/systemd/rclone-managed-bisync.timer
[Unit]
Description=Run managed rclone bisync jobs every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
